# docker/cron.d/sslrenew
0 0 * * 1 root certbot renew >> /var/log/certbot-renew.log 2>&1
